# server/crontabs/lt-crons
# * * * * * cd /usr/www/src/laravel && php artisan schedule:run >> /dev/null 2>&1
* * * * * cd /usr/www/src/mclassvn3 && php artisan schedule:run >> /dev/null 2>&1